# services/socialfy/AgenticOSKevsAcademy/com.agenticos.dailyanalytics.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agenticos.dailyanalytics</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/newnewmacbookacbook/AGENTIC OS OFFICAL/agenticos/implementation/instantly_analytics.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>20</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>WorkingDirectory</key>
    <string>/Users/newnewmacbookacbook/AGENTIC OS OFFICAL/agenticos</string>
    <key>StandardOutPath</key>
    <string>/Users/newnewmacbookacbook/AGENTIC OS OFFICAL/agenticos/logs/daily_analytics.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/newnewmacbookacbook/AGENTIC OS OFFICAL/agenticos/logs/daily_analytics_error.log</string>
</dict>
</plist>