# services/socialfy/AgenticOSKevsAcademy/com.socialfy.campaignworker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.socialfy.campaignworker</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/marcosdaniels/Projects/mottivme/1. ai-factory-mottivme-sales/4. socialfy-platform/AgenticOSKevsAcademy/implementation/campaign_worker.py</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/marcosdaniels/Projects/mottivme/1. ai-factory-mottivme-sales/4. socialfy-platform/AgenticOSKevsAcademy/implementation</string>
    
    <key>StartInterval</key>
    <integer>60</integer>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <false/>
    
    <key>StandardOutPath</key>
    <string>/Users/marcosdaniels/Library/Logs/socialfy/campaign_worker.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/marcosdaniels/Library/Logs/socialfy/campaign_worker_error.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>SUPABASE_URL</key>
        <string>https://bfumywvwubvernvhjehk.supabase.co</string>
        <key>SUPABASE_SERVICE_KEY</key>
        <string>eyJhbGciOiJIUzI1NiIsInR5cCI6IkpXVCJ9.eyJpc3MiOiJzdXBhYmFzZSIsInJlZiI6ImJmdW15d3Z3dWJ2ZXJudmhqZWhrIiwicm9sZSI6InNlcnZpY2Vfcm9sZSIsImlhdCI6MTc1MTQwMzc5OSwiZXhwIjoyMDY2OTc5Nzk5fQ.fdTsdGlSqemXzrXEU4ov1SUpeDn_3bSjOingqkSAWQE</string>
    </dict>
</dict>
</plist>
